# install/kingshot-bot.service
[Unit]
Description=Kingshot Discord Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/kingshot-bot
# --no-update: auto-update would overwrite this fork's custom cogs with upstream kingshot-project files
ExecStart=/home/pi/kingshot-bot/bot_venv/bin/python main.py --no-update
Restart=on-failure
RestartSec=10
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
